# config/systemd/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=%h/.local/share/.dropbox-dist/dropboxd
Restart=on-failure

[Install]
WantedBy=default.target
